# public/scripts/gwi-node.service
[Unit]
Description=GWI Node Deploy Agent
Documentation=https://github.com/GetwithitMan/gwi-pos
After=docker.service
Requires=docker.service
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
ExecStart=/opt/gwi-pos/gwi-node.sh watch
WorkingDirectory=/opt/gwi-pos
Restart=always
RestartSec=5
User=root
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gwi-node

# Hardening
ProtectSystem=full
NoNewPrivileges=no
PrivateTmp=true

[Install]
WantedBy=multi-user.target
