# installer/roles/app_services/templates/gwi-converge.service.j2
[Unit]
Description=GWI POS Venue Convergence Agent
After=docker.service
Requires=docker.service

[Service]
Type=simple
ExecStart={{ gwi_base_dir | default('/opt/gwi-pos') }}/gwi-node.sh converge-loop 300
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gwi-converge

[Install]
WantedBy=multi-user.target
