# installer/roles/app_services/templates/gwi-converge.timer.j2
# DEPRECATED: gwi-converge now uses a long-running service (converge-loop),
# not a timer. This file is kept only to prevent Ansible errors during
# transition. The timer is not enabled — only gwi-converge.service is used.
# See gwi-node.sh converge-loop subcommand.
[Unit]
Description=GWI Converge Timer (DEPRECATED — use gwi-converge.service directly)

[Timer]
OnBootSec=999h

[Install]
WantedBy=timers.target
